# config/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /home/docker/code/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx

stdout_logfile="~/stdout_logfile.log"
stderr_logfile="~/stderr_logfile,log"
autostart=true
autorestart=true
startsecs=10